# build/syncthing.conf
[program:syncthing]
user = nobody
command = /usr/bin/syncthing -config="/config/syncthing/config" -data="/config/syncthing/data" -gui-address="http://0.0.0.0:8384"
umask = 000
